# deploy/systemd/kassinao-health-watch.service
[Unit]
Description=Recover unhealthy Kassinão components without exposing the Docker socket
After=docker.service kassinao-docker-egress.service
Requires=docker.service kassinao-docker-egress.service
ConditionPathExists=/var/run/docker.sock
OnFailure=kassinao-egress-fail-closed.service
OnFailureJobMode=replace-irreversibly

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/kassinao-health-watch
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/lock/kassinao
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK
RestrictNamespaces=true
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
LockPersonality=true
MemoryDenyWriteExecute=true
