# deploy/systemd/kassinao-health-watch.timer
[Unit]
Description=Check Kassinão Docker health every two minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=15s
Persistent=true
Unit=kassinao-health-watch.service

[Install]
WantedBy=timers.target
