# deploy/systemd/kassinao-health-watch.service
[Unit]
Description=Restart Kassinão only when its Docker healthcheck is unhealthy
After=docker.service
Requires=docker.service
ConditionPathExists=/var/run/docker.sock

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/kassinao-health-watch
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/lock
RestrictAddressFamilies=AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true
