# deploy/systemd/docker.service.d/kassinao-egress.conf
[Unit]
Requires=systemd-tmpfiles-setup.service
After=systemd-tmpfiles-setup.service
Wants=kassinao-docker-egress.service

[Service]
ExecStartPre=/usr/local/sbin/kassinao-verify-storage-encryption
ExecStartPre=/usr/local/sbin/kassinao-harden-docker-egress --offline-preload
